# devservices/programs.conf
[program:devserver]
command=snuba devserver
autostart=false
